# infra/systemd/quant-core.service
[Unit]
Description=Quant-Core stack (docker compose)
Wants=docker.service
After=docker.service

[Service]
Type=simple
WorkingDirectory={{ROOT_DIR}}
Environment="ENV_FILE={{ENV_FILE}}"
User=root
Group=root
ExecStart=/usr/bin/env bash -c '{{ROOT_DIR}}/infra/scripts/prod-deploy.sh --no-pull && exec sleep infinity'
ExecReload=/usr/bin/env bash {{ROOT_DIR}}/infra/scripts/prod-deploy.sh --no-pull
ExecStop=/usr/bin/env docker compose -f {{ROOT_DIR}}/infra/docker-compose.yml --env-file {{ENV_FILE}} down
Restart=always
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
